# docker-createrepo/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

#[inet_http_server]
#port = 9001
#username = user
#password = pass 

[program:nginx_1.9.9]
command=/bin/bash -c "exec /usr/local/nginx/sbin/nginx -g 'daemon off;'"
autostart=true
autorestart=false
startretries=0
stdout_events_enabled=true
stderr_events_enabled=true

[program:nginx_1.9.9_conf]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /data/wwwroot/; do createrepo -update /data/wwwroot/; done'
autostart=true
autorestart=false
startretries=0
stdout_events_enabled=true
stderr_events_enabled=true
